# depecher/systemd/org.blacksailer.depecher.service
[Unit]
Description=native application for Telegram based on tdlib
Requires=lipstick.service dbus.service
After=lipstick.service dbus.service

[Service]
Type=dbus
ExecStart=/usr/bin/invoker  -s --type=silica-qt5 /usr/bin/depecher
BusName=org.blacksailer.depecher

[Install]
WantedBy=post-user-session.target
